# tempest-mqtt.service
[Unit]
Description=WeatherFlow Tempest to MQTT Daemon
After=network.target

[Service]
Type=simple
# Config file: ~/.config/tempest-mqtt.env
# MQTT_USERNAME=myuser
# MQTT_PASSWORD=secret
# TOPIC_PREFIX=home          # → topics: climate/home/...
EnvironmentFile=-%h/.config/tempest-mqtt.env
ExecStart=%h/tempest-mqtt/tempest-mqtt \
    -broker ${MQTT_BROKER_URL} \
    -client-id ${MQTT_CLIENT_ID} \
    -username ${MQTT_USERNAME} \
    -password ${MQTT_PASSWORD} \
    -topic-prefix ${TOPIC_PREFIX} \
    -udp-port ${TEMPEST_UDP_PORT} \
    -log-level ${LOG_LEVEL}
Restart=on-failure
RestartSec=5
TimeoutStopSec=10

[Install]
WantedBy=default.target
